# g3-access-dashboard/deploy/supervisor/g3-access-worker.conf
; Supervisor config for G3 Access queue workers.
;
; Install: copy to /etc/supervisor/conf.d/g3-access-worker.conf on the EC2 instance,
; adjust paths/user, then:
;   sudo supervisorctl reread
;   sudo supervisorctl update
;   sudo supervisorctl start g3-access-worker:*
;
; Restart after deploys:
;   sudo supervisorctl restart g3-access-worker:*

[program:g3-access-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/g3-access-dashboard/artisan queue:work --queue=default --sleep=3 --tries=1 --max-time=3600 --max-jobs=100 --timeout=180
autostart=true
autorestart=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/log/supervisor/g3-access-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stopwaitsecs=200
